# vendor/cookbooks/scalient/templates/ubuntu/chef-client.conf.erb
# chef-client - Chef Client
#
# Chef Client runs the Chef client on system startup
#
#  ____________________________________
# / This file is managed by Chef. Your \
# \ changes will be overwritten.       /
#  ------------------------------------
#         \   ^__^
#          \  (oo)\_______
#             (__)\       )\/\
#                 ||----w |
#                 ||     ||
#

description "Chef Client"

<% if Pathname.new("/etc/init.d/chef-server").file? %>
start on (filesystem and started chef-server)
<% else %>
start on filesystem
<% end %>

stop on runlevel [!2345]

env RBENV_VERSION=<%= @rbenv_version.shellescape %>

pre-start script
    test -x "`rbenv which chef-client`" || { stop; exit 0; }
end script

script
    prefix=<%= @prefix.shellescape %>

    exec rbenv exec chef-client \
        -c "${prefix}/etc/chef-client/client.rb" \
        -L "${prefix}/var/log/chef-client/client.log" \
        -o "recipe[scalient::aws]"
end script
