# vendor/cookbooks/scalient/templates/ubuntu/chromium-browser.service.erb
#  ____________________________________
# / This file is managed by Chef. Your \
# \ changes will be overwritten.       /
#  ------------------------------------
#         \   ^__^
#          \  (oo)\_______
#             (__)\       )\/\
#                 ||----w |
#                 ||     ||
#

[Unit]
Description=The Chrome headless browser is used to render (and print) web pages in a consistent manner

[Service]
Type=simple
WorkingDirectory=<%= @app_root.shellescape %>
User=<%= @original_user.shellescape %>
ExecStart=/usr/bin/chromium-browser --headless --remote-debugging-port=9222 https://chromium.org/

[Install]
WantedBy=multi-user.target
