# vendor/cookbooks/scalient/templates/ubuntu/unicorn.service.erb
#  ____________________________________
# / This file is managed by Chef. Your \
# \ changes will be overwritten.       /
#  ------------------------------------
#         \   ^__^
#          \  (oo)\_______
#             (__)\       )\/\
#                 ||----w |
#                 ||     ||
#

[Unit]
Description=Unicorn is an HTTP server for Rack applications.
Requires=postgresql.service

[Service]
Type=simple
WorkingDirectory=<%= @app_root.shellescape %>
Environment=RBENV_VERSION=<%= @rbenv_version.shellescape %>
User=<%= @original_user.shellescape %>
ExecStart=/usr/bin/rbenv exec bundle exec unicorn_rails -E production -c config/unicorn.rb

[Install]
WantedBy=multi-user.target
