# vendor/cookbooks/scalient/templates/ubuntu/unicorn.conf.erb
# unicorn - Unicorn HTTP server
#
# Unicorn is an HTTP server for Rack applications.
#
#  ____________________________________
# / This file is managed by Chef. Your \
# \ changes will be overwritten.       /
#  ------------------------------------
#         \   ^__^
#          \  (oo)\_______
#             (__)\       )\/\
#                 ||----w |
#                 ||     ||
#

description "Unicorn"

start on filesystem
stop on runlevel [!2345]

chdir <%= Pathname.new(@app_root).to_s.shellescape %>
env RBENV_VERSION=<%= @rbenv_version.shellescape %>

pre-start script
    test -f <%= Pathname.new(@app_root).join("Gemfile").to_s.shellescape %> || { stop; exit 0; }
end script

script
    exec sudo -u <%= @original_user.shellescape %> -- rbenv exec bundle exec unicorn_rails \
        -E production \
        -c config/unicorn.rb \
end script
